# pi_pad.service
[Unit]
Description=Pi Pad - Web Touchpad Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/YOUR_USERNAME/pi_pad
ExecStart=/home/YOUR_USERNAME/pi_pad/venv/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000
ExecStop=/home/YOUR_USERNAME/pi_pad/setup_gadget.sh --clean
Restart=always

[Install]
WantedBy=multi-user.target


